# roles/pinger/files/pinger.service
[Unit]
After=network-online.target

[Service]
ExecStart=/root/pinger

[Install]
WantedBy=multi-user.target